# web_scraper/crontab
00 08 * * * cd /workspace/web_scraper && TZ=America/Toronto /usr/bin/python3 ./utils/weather.py >> /workspace/logs/cron.log 2>&1 
58 08 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Altitude_Gatineau >> /workspace/logs/cron.log 2>&1  
18 11 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Altitude_Gatineau >> /workspace/logs/cron.log 2>&1  
38 13 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Altitude_Gatineau >> /workspace/logs/cron.log 2>&1  
58 15 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Altitude_Gatineau >> /workspace/logs/cron.log 2>&1  
18 18 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Altitude_Gatineau >> /workspace/logs/cron.log 2>&1  
38 20 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Altitude_Gatineau >> /workspace/logs/cron.log 2>&1  
58 08 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Altitude_Kanata >> /workspace/logs/cron.log 2>&1  
18 11 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Altitude_Kanata >> /workspace/logs/cron.log 2>&1  
38 13 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Altitude_Kanata >> /workspace/logs/cron.log 2>&1  
58 15 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Altitude_Kanata >> /workspace/logs/cron.log 2>&1  
18 18 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Altitude_Kanata >> /workspace/logs/cron.log 2>&1  
38 20 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Altitude_Kanata >> /workspace/logs/cron.log 2>&1  
28 10 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Coyote_Rock_Gym >> /workspace/logs/cron.log 2>&1  
43 12 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Coyote_Rock_Gym >> /workspace/logs/cron.log 2>&1  
58 14 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Coyote_Rock_Gym >> /workspace/logs/cron.log 2>&1  
13 17 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Coyote_Rock_Gym >> /workspace/logs/cron.log 2>&1  
28 19 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Coyote_Rock_Gym >> /workspace/logs/cron.log 2>&1  
43 21 * * 1-5 cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py -l Coyote_Rock_Gym >> /workspace/logs/cron.log 2>&1  
