# web_scraper/crontab
18 09 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py >> /workspace/logs/cron.log 2>&1  
38 11 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py >> /workspace/logs/cron.log 2>&1  
58 13 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py >> /workspace/logs/cron.log 2>&1  
18 16 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py >> /workspace/logs/cron.log 2>&1  
38 18 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py >> /workspace/logs/cron.log 2>&1  
58 20 * * * cd /workspace/web_scraper && DOCKER_SCRAPER=true TZ=America/Toronto /usr/bin/python3 ./bookings.py >> /workspace/logs/cron.log 2>&1  
